# deploy/systemd/ovcp-openvpn.service
[Unit]
Description=OpenVPN server (OVCP)
# bound to ovcp.service: started, stopped and restarted together with it
BindsTo=ovcp.service
After=ovcp.service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/sbin/openvpn --config /var/lib/ovcp/server.conf
# always, not on-failure: ovcp's "restart" is a SIGTERM (clean exit 0)
# through the mgmt socket; systemd must respawn openvpn fresh as root.
Restart=always
RestartSec=2
AmbientCapabilities=CAP_NET_ADMIN
DeviceAllow=/dev/net/tun rw
RuntimeDirectory=ovcp
RuntimeDirectoryMode=0750
Group=ovcp

[Install]
WantedBy=multi-user.target
