# deploy/systemd/ovcp.service
[Unit]
Description=OVCP admin UI / control plane
# pulls openvpn in on start; the reverse BindsTo restarts it with us
Wants=ovcp-openvpn.service
Before=ovcp-openvpn.service

[Service]
User=ovcp
Group=ovcp
ExecStart=/usr/local/bin/ovcp serve
Environment=OVCP_DATA=/var/lib/ovcp
Environment=OVCP_PLATFORM=systemd
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ovcp

[Install]
WantedBy=multi-user.target
