# deploy/systemd/ovcp.service
[Unit]
Description=OVCP admin UI / control plane
Wants=openvpn-ovcp.service
After=openvpn-ovcp.service

[Service]
Type=simple
User=ovcp
Group=ovcp
ExecStart=/usr/local/bin/ovcp serve
Environment=OVCP_DATA=/var/lib/ovcp
Environment=OVCP_PLATFORM=systemd
Environment=OVCP_SYSTEMD_UNIT=openvpn-ovcp
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ovcp
ProtectHome=true

[Install]
WantedBy=multi-user.target
